# scripts/playground/playground-monitor.service
[Unit]
Description=Dakera Playground Health Monitor — DAK-6745
After=network.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/dakera/telegram.env
ExecStart=/usr/local/bin/playground-monitor.sh
User=root
StandardOutput=journal
StandardError=journal
