# scripts/runner/runner-health-monitor.service
[Unit]
Description=Dakera Runner Health Monitor — DAK-5764
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/etc/dakera/telegram.env
ExecStart=/usr/local/bin/runner-health-monitor.sh
User=root
StandardOutput=journal
StandardError=journal
