# scripts/runner/runner-health-monitor.timer
[Unit]
Description=Dakera Runner Health Monitor Timer — DAK-5764

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
